# deploy/systemd/roguelike-lemonade-stand.service
[Unit]
Description=Roguelike Lemonade Stand multiplayer server
After=network.target

[Service]
Type=simple
# Replace these with the user and group that own the app checkout.
User=deploy
Group=deploy
WorkingDirectory=/var/www/roguelike-lemonade-stand
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/npm run server:start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
